# stow/.config/systemd/user/whisper-ptt.service
[Unit]
Description=Whisper Push-to-Talk Daemon
After=graphical-session.target pipewire.service
Requires=graphical-session.target

[Service]
Type=simple
Environment=CUDA_VISIBLE_DEVICES=0
Environment=PATH=%h/.local/bin:%h/.venv/whisper/bin:/run/current-system/sw/bin:/usr/bin
ExecStart=%h/.venv/whisper/bin/python %h/.config/whisper-ptt/whisper-daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
